# Raspberrypi/server/server.service
[Unit]
Description=Polar H10 server (ESP32 → MQTT → SQLite) for the Pi
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__DIR__
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=__DIR__/mqtt.env
ExecStart=__DIR__/.venv/bin/python3 __DIR__/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
