# Raspberrypi/dashboard/dashboard.service
[Unit]
Description=Polar H10 web dashboard (reads hr_data.db, serves live view)
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__DIR__
Environment=PYTHONUNBUFFERED=1
Environment=PORT=8000
# HR_DB defaults to ../server/hr_data.db — set it here to override.
ExecStart=__DIR__/.venv/bin/python3 __DIR__/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
